# conf.d/3-websockify.conf
[program:websockify]
command=bash -c 'websockify --cert /etc/ssl/novnc.cert --key /etc/ssl/novnc.key -D --web=/usr/share/novnc/ $PORT localhost:5900'
autostart=true
autorestart=true
stderr_logfile=/var/log/websockify.stderr.log
priority=997
