# conf.d/4-fluxbox.conf
[program:Fluxbox]
command=fluxbox
autostart=true
autorestart=true
priority=996
stderr_logfile=/var/log/fluxbox.err.log