# assets/configs/1-xvfb.conf
[program:Xvfb]
command=bash -c 'Xvfb $DISPLAY -screen 0 $VNC_RESOLUTION"x24" -nolisten unix'
stderr_logfile=/app/.vubuntu/assets/logs/Xvfb.stderr.log
priority=1
