# assets/configs/4-fluxbox.conf
[program:Fluxbox]
command=fluxbox
exitcodes=1
priority=4
stderr_logfile=/app/.vubuntu/assets/logs/fluxbox.err.log
