# assets/configs/5-port-forward.conf
[program:portforward]
command=bash -c 'socat TCP-LISTEN:$PORT,fork TCP:127.0.0.1:4040'
stderr_logfile=/var/log/socat.err.log
autorestart=unexpected
exitcodes=0
startseconds=0
priority=1
