# assets/configs/5-selfping.conf
[program:self-ping]
command=python3 /app/.vubuntu/assets/scripts/self-ping.py
stderr_logfile=/var/log/selfping.err.log
stdout_logfile=/var/log/selfping.out.log
priority=995
