# extras/4-mate.conf
[program:Mate]
command=mate-session
autostart=true
autorestart=true
stderr_logfile:/var/log/mate.err.log
priority=996