# extras/4-xfce.conf
[program:Xfce]
command=startxfce4
autostart=true
autorestart=true
stderr_logfile:/var/log/xfce.err.log
priority=996